# ops/systemd/user/clawtext-daily-cluster-rebuild.timer
[Unit]
Description=Run ClawText daily cluster rebuild at 02:00 UTC

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true
Unit=clawtext-daily-cluster-rebuild.service

[Install]
WantedBy=timers.target
